# docker/src/supervisor/conf.d/system.conf
[program:system]
command=bash system-init.sh
user=root
directory=/dockerstartup
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/var/log/supervisor/system-init.log
redirect_stderr=true
environment=
   APPUSER_PASSWORD=%(ENV_APPUSER_PASSWORD)s,
   APT_INSTALL=%(ENV_APT_INSTALL)s

   